# papercraft/etc/systemd/system/container-getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noclear --autologin sakamoto --keep-baud pts/%I 115200,38400,9600 $TERM